# Automation/systemd/seer-move-oldest.service
[Unit]
Description=SEER mover: when ≥ threshold, move oldest closed pcap to dest (Req3)
After=network-online.target

[Service]
Type=oneshot
User=seer
Group=seer
ExecStart=/usr/bin/env python3 /home/seer/code/SEER-Sensor/Automation/SEER/move_oldest.py

