# Automation/systemd/seer-capture@.service
[Unit]
Description=SEER tcpdump capture on %i (Req1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root so tcpdump can open the interface; it will drop to user 'seer' via -Z
ExecStart=/usr/local/bin/seer-capture.sh %i
KillSignal=SIGINT
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

